# scripts/automation/com.hxz.news-reader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hxz.news-reader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/admin/Downloads/hxz-news-reader/scripts/automation/run_local.py</string>
    </array>

    <!-- 每天本地时间 06:30 跑 (GitHub Actions cron 6:00 完成后 30 分钟) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- 错过时间补跑 (合盖醒来等) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- launchd 不读 ~/.zshrc, 必须在这里把 X cookie 注入子进程 -->
    <!-- 注意: 这个 plist 含敏感 cookie, 一定不能 git commit -->
    <!-- 真实值在用户 home 副本里, 仓库这版字段是占位符 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>TWITTER_AUTH_TOKEN</key>
        <string>PLACEHOLDER_REPLACE_AT_INSTALL</string>
        <key>TWITTER_CT0</key>
        <string>PLACEHOLDER_REPLACE_AT_INSTALL</string>
        <key>NO_PROXY</key>
        <string>127.0.0.1,localhost</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/admin/Library/Logs/hxz-news-reader.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/Library/Logs/hxz-news-reader.stderr.log</string>
</dict>
</plist>
